# .taskfiles/Machine/files/artemis/cloudflare-ddns.crontab
*/10 * * * * root /usr/local/bin/cloudflare-ddns.sh
